# root/etc/services.d/chowdown/run
#!/usr/bin/with-contenv bash

cd /app/chowdown || exit

JEKYLL_ENV=production

exec \
	s6-setuidgid abc /usr/local/bin/jekyll serve --host 0.0.0.0 --baseurl "$BASEURL"